# debian/mailrobots-ham-learn.service
[Unit]
Description=ham learn

[Service]
Type=oneshot
User=root
ExecStart=/usr/lib/mailrobots/ham-learn
