# debian/mailrobots-spam-purge.service
[Unit]
Description=spam/trash purge

[Service]
Type=oneshot
User=root
ExecStart=/usr/lib/mailrobots/spam-purge
